# packaging/deb/prerm
#!/bin/sh
set -e

systemctl stop openusbd.service 2>/dev/null || true
systemctl disable openusbd.service 2>/dev/null || true
